# dockyard/supervisord.d/sshd.ini
; Applications
[program:sshd]
;directory=/usr/local/
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startsecs=10
redirect_stderr=true
priority=910
